# systemd/org-commit.timer
[Unit]
Description=Commit ~/org/ every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=org-commit.service

[Install]
WantedBy=timers.target
